# build-rootfs/files/etc/inittab
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# Gettys for every tty Podroid might use. podroid-getty selects which to run:
# - hvc0/ttyS0: gated by the `podroid.tty=` kernel cmdline marker, so the
#   inactive backend's tty parks itself; keeps the squashfs backend-agnostic.
# - hvc2/hvc3: extra multi-tab terminals exposed by the QEMU backend's
#   virtio-console bus. The wrapper parks them when the device node is missing
#   (AVF backend doesn't wire these), so they're harmless there.
hvc0::respawn:/usr/local/bin/podroid-getty hvc0
ttyS0::respawn:/usr/local/bin/podroid-getty ttyS0
hvc2::respawn:/usr/local/bin/podroid-getty hvc2
hvc3::respawn:/usr/local/bin/podroid-getty hvc3
::shutdown:/sbin/openrc shutdown
